# profile.d/color.sh
if [ -n "$IS_COLOR" ]; then
  export COLOR_RESET='0m'
  export COLOR_BLACK='0;30m'
  export COLOR_RED='0;31m'
  export COLOR_GREEN='0;32m'
  export COLOR_YELLOW='0;33m'
  export COLOR_BLUE='01;34m'
  export COLOR_PURPLE='0;35m'
  export COLOR_CYAN='0;36m'
  export COLOR_WHITE='0;37m'
fi
