# profile.d/docker.sh
alias db='docker build --rm'
alias dcb='docker-compose build'
alias dcd='docker-compose down'
alias dce='docker-compose exec'
alias dcm='docker-compose exec app python manage.py'
alias dcr='docker-compose run --rm'
alias dcu='docker-compose up'
alias dcub='docker-compose up --build'
alias de='docker exec'
alias dei='docker exec -it'
alias dr='docker run --rm'
alias dri='docker run --rm -it'
alias drm='docker container prune -f; docker image prune -f; docker volume prune -f; docker network prune -f;'
