# profile.d/npm.sh
alias npminit='echo "{\n  \"private\": true\n}" > package.json'
