# profile.d/platform.sh
[ -r "$HOME/.profile-$PLATFORM" ] && . "$HOME/.profile-$PLATFORM"
